# configs/macos/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": false,
    "log_result_events": true,
    "schedule_splay_percent": 10,
    "pidfile": "/var/osquery/osquery.pidfile",
    "events_expiry": 3600,
    "database_path": "/var/osquery/osquery.db",
    "verbose": false,
    "worker_threads": 2,
    "enable_monitor": true,
    "disable_events": false,
    "disable_audit": false,
    "audit_allow_config": true,
    "audit_allow_sockets": true,
    "host_identifier": "hostname",
    "enable_syslog": true,
    "utc": true
  },
  
  "schedule": {
    "network_connections_established": {
      "query": "SELECT p.name AS process_name, p.path AS process_path, p.pid, p.cmdline, u.username, pos.local_address, pos.local_port, pos.remote_address, pos.remote_port, pos.protocol, pos.state FROM processes p JOIN process_open_sockets pos USING (pid) LEFT JOIN users u ON p.uid = u.uid WHERE pos.remote_port != 0 AND pos.state = 'ESTABLISHED';",
      "interval": 300,
      "description": "Active established network connections with process info",
      "snapshot": true
    },
    
    "network_connections_all": {
      "query": "SELECT p.name AS process_name, p.path AS process_path, p.pid, pos.local_address, pos.local_port, pos.remote_address, pos.remote_port, pos.protocol, pos.state FROM processes p JOIN process_open_sockets pos USING (pid) WHERE pos.remote_port != 0;",
      "interval": 60,
      "description": "All network connections",
      "snapshot": true
    },
    
    "listening_ports": {
      "query": "SELECT p.name AS process_name, p.path AS process_path, p.pid, p.cmdline, lp.port, lp.address, lp.protocol FROM processes p JOIN listening_ports lp USING (pid);",
      "interval": 600,
      "description": "Listening ports with process attribution",
      "snapshot": true
    },
    
    "listening_ports_anomaly": {
      "query": "SELECT p.name, p.path, lp.port, lp.address, lp.protocol FROM processes p JOIN listening_ports lp USING (pid) WHERE lp.port NOT IN (22, 80, 443, 53, 3306, 5432, 6379, 8080) AND lp.address IN ('0.0.0.0', '::');",
      "interval": 300,
      "description": "Non-standard listening ports",
      "snapshot": true
    },
    
    "dns_resolvers": {
      "query": "SELECT * FROM dns_resolvers;",
      "interval": 3600,
      "description": "DNS resolver configuration",
      "snapshot": true
    },
    
    "network_interfaces": {
      "query": "SELECT interface, address, mask, broadcast, type FROM interface_addresses WHERE type = 'IPv4' OR type = 'IPv6';",
      "interval": 3600,
      "description": "Network interface configuration",
      "snapshot": true
    },
    
    "routes": {
      "query": "SELECT destination, gateway, interface, type FROM routes WHERE destination != '::1' AND destination != '127.0.0.1';",
      "interval": 3600,
      "description": "Routing table",
      "snapshot": true
    },
    
    "arp_cache": {
      "query": "SELECT address, mac, interface FROM arp_cache;",
      "interval": 1800,
      "description": "ARP cache entries",
      "snapshot": true
    },
    
    "external_connections": {
      "query": "SELECT DISTINCT p.name, p.path, pos.remote_address, pos.remote_port FROM processes p JOIN process_open_sockets pos USING (pid) WHERE pos.remote_address NOT LIKE '127.%' AND pos.remote_address NOT LIKE '10.%' AND pos.remote_address NOT LIKE '192.168.%' AND pos.remote_address NOT LIKE '172.16.%' AND pos.remote_address NOT LIKE '172.17.%' AND pos.remote_address NOT LIKE '172.18.%' AND pos.remote_address NOT LIKE '172.19.%' AND pos.remote_address NOT LIKE '172.2%' AND pos.remote_address NOT LIKE '172.30.%' AND pos.remote_address NOT LIKE '172.31.%' AND pos.remote_port != 0;",
      "interval": 120,
      "description": "Connections to external (non-RFC1918) addresses",
      "snapshot": true
    },
    
    "suspicious_process_network": {
      "query": "SELECT p.name, p.path, p.cmdline, p.pid, pos.remote_address, pos.remote_port FROM processes p JOIN process_open_sockets pos USING (pid) WHERE (p.name IN ('curl', 'wget', 'nc', 'ncat', 'netcat', 'socat', 'python', 'python3', 'perl', 'ruby', 'php') OR p.path LIKE '/tmp/%' OR p.path LIKE '/var/tmp/%') AND pos.remote_port != 0;",
      "interval": 60,
      "description": "Network activity from suspicious processes",
      "snapshot": true
    },
    
    "process_events": {
      "query": "SELECT pid, path, cmdline, uid, gid, parent, time FROM process_events;",
      "interval": 60,
      "description": "Process execution events",
      "removed": false
    },
    
    "user_ssh_keys": {
      "query": "SELECT * FROM user_ssh_keys;",
      "interval": 86400,
      "description": "User SSH keys",
      "snapshot": true
    },
    
    "authorized_keys": {
      "query": "SELECT * FROM authorized_keys;",
      "interval": 3600,
      "description": "SSH authorized keys",
      "snapshot": true
    },
    
    "crontab": {
      "query": "SELECT * FROM crontab;",
      "interval": 3600,
      "description": "Scheduled cron jobs",
      "snapshot": true
    },
    
    "launchd_agents": {
      "query": "SELECT name, path, program, program_arguments, run_at_load FROM launchd WHERE run_at_load = 1;",
      "interval": 3600,
      "description": "Persistent launchd agents",
      "snapshot": true,
      "platform": "darwin"
    },
    
    "browser_plugins": {
      "query": "SELECT * FROM browser_plugins;",
      "interval": 86400,
      "description": "Browser plugins",
      "snapshot": true,
      "platform": "darwin"
    }
  },
  
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT hostname AS hostname FROM system_info;",
      "SELECT version AS osquery_version FROM osquery_info;"
    ]
  },
  
  "packs": {
  },
  
  "file_paths": {
    "etc": [
      "/etc/%%"
    ],
    "homes": [
      "/Users/%/.ssh/%%",
      "/Users/%/.bash_history"
    ],
    "tmp": [
      "/tmp/%%"
    ]
  }
}
